# system-install/profile/iso/configs/arch-iso-net.service
[Unit]
Description=iso script
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=greetd

[Install]
WantedBy=multi-user.target
